# config/plasmaConfSaver/default/autostart/io.github.vikdevelop.SaveDesktop.Backup.desktop
[Desktop Entry]
Name=SaveDesktop (Periodic backups)
Type=Application
Exec=flatpak run io.github.vikdevelop.SaveDesktop --background